# scripts/tester.service
[Unit]
Description=Tester Service
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=tester
ExecStart=/usr/bin/nice -19 /usr/local/bin/tester
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Load environment variables from configuration file
EnvironmentFile=-/etc/tester/tester.env

# XDG Base Directory Specification (optional, will use defaults if not set)
Environment="XDG_DATA_DIRS=/usr/local/share:/usr/share"
Environment="XDG_CONFIG_DIRS=/etc/xdg"

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
